# timers/change-wallpaper.service
[Unit]
Description=Change wallpaper at regular intervals

[Service]
ExecStart=/home/will/.changebg
