# timers/change-wallpaper.timer
[Unit]
Description=Change wallpaper every 15 mins

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
