# timers/bgchange/bgchange.service
[Unit]
Description=Change background image periodically

[Service]
ExecStart=/home/will/.bgchange
