# timers/bgchange/bgchange.timer
[Unit]
Description=Timer for background change

[Timer]
OnActiveSec=5sec
OnUnitActiveSec=15min
